# FILES/HID.service
[Unit]
Description=volumio-hid Service
Documentation=https://github.com/edmw/volumio-hid
After=volumio.target

[Service]
ExecStart=/home/volumio/.pyenv/volumio-hid/bin/python /home/volumio/volumio-hid/HID.py
User=volumio
Group=volumio

[Install]
WantedBy=multi-user.target
